# devops/systemd/crypto-thing.service
[Unit]
Description=Crypto Thing Tool
After=network.target
Wants=network.target

[Service]
Type=simple
User=crypto
Group=crypto
WorkingDirectory=/opt/crypto-thing
EnvironmentFile=/opt/crypto-thing/.env

# Start the daemon with websocket server
ExecStart=/opt/crypto-thing/cryptool daemon --port 40000

# Restart policy
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/crypto-thing /etc/crypto-thing

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crypto-thing

[Install]
WantedBy=multi-user.target
